# core/supervisord/conf.d/api.conf
[program:api]
command=uvicorn main:app --host 0.0.0.0 --port 8091
directory=/opt/sandbox/api
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
environment=DISPLAY=":1"
