# core/supervisord/conf.d/chromium.conf
[program:chromium]
; Playwright-bundled Chromium — works on arm64 without snap.
; Path is written to /opt/chromium-path at build time.
command=bash -c "rm -f /root/.config/chrome-sandbox/Singleton* && exec $(cat /opt/chromium-path) --no-sandbox --disable-dev-shm-usage --disable-gpu --use-gl=swiftshader --remote-debugging-port=9222 --remote-debugging-address=0.0.0.0 --remote-allow-origins=* --user-data-dir=/root/.config/chrome-sandbox --window-size=1280,800 --window-position=0,0 --no-first-run --no-default-browser-check --disable-background-networking --disable-default-apps --disable-extensions-except=/opt/sandbox/stealth-extension --load-extension=/opt/sandbox/stealth-extension --disable-sync --metrics-recording-only --password-store=basic --disable-blink-features=AutomationControlled about:blank"
autostart=true
autorestart=true
priority=40
startsecs=5
startretries=3
user=root
environment=DISPLAY=":1",HOME="/root",DBUS_SESSION_BUS_ADDRESS="autolaunch:"
stdout_logfile=/var/log/supervisor/chromium.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/chromium-err.log
stderr_logfile_maxbytes=5MB
