# core/supervisord/conf.d/mcp.conf
[program:mcp]
command=uvicorn server:create_app --factory --host 0.0.0.0 --port 8079
directory=/opt/sandbox/mcp_server
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp.log
stderr_logfile=/var/log/supervisor/mcp.log
environment=DISPLAY=":1",CDP_URL="http://localhost:9222",ADB_HOST="localhost"
