# core/supervisord/conf.d/nginx.conf
[program:nginx]
; nginx reverse proxy / dashboard on port 8080
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
startsecs=2
startretries=5
user=root
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/nginx-err.log
stderr_logfile_maxbytes=5MB
