# core/supervisord/conf.d/novnc.conf
[program:novnc]
; websockify bridges VNC on localhost:5900 to WebSocket on 0.0.0.0:6080
; noVNC static files are served from /opt/novnc
command=/usr/bin/websockify
    --web=/opt/novnc
    --heartbeat=30
    0.0.0.0:6080
    localhost:5900
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=5
user=root
stdout_logfile=/var/log/supervisor/novnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/novnc-err.log
stderr_logfile_maxbytes=5MB
