# core/supervisord/conf.d/vnc.conf
[program:vnc]
; TigerVNC standalone server on display :1, geometry 1280x800, depth 24
command=/usr/bin/tigervncserver :1
    -geometry 1280x800
    -depth 24
    -rfbport 5900
    -localhost no
    -SecurityTypes VncAuth
    -xstartup /etc/xdg/xfce4/xinitrc
autostart=true
autorestart=true
priority=10
startsecs=3
startretries=5
user=root
environment=HOME="/root",USER="root"
stdout_logfile=/var/log/supervisor/vnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/vnc-err.log
stderr_logfile_maxbytes=5MB
