# core/supervisord/conf.d/vnc.conf
[program:xtigervnc]
; Xtigervnc X server on display :1, no auth (local sandbox)
command=/usr/bin/Xtigervnc :1 -geometry 1280x800 -depth 24 -rfbport 5900 -localhost no -SecurityTypes None
autostart=true
autorestart=true
priority=10
startsecs=3
startretries=5
user=root
environment=HOME="/root",USER="root"
stdout_logfile=/var/log/supervisor/vnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/vnc-err.log
stderr_logfile_maxbytes=5MB

[program:xfce]
; XFCE4 desktop session on display :1, started after xtigervnc
command=/usr/bin/xfce4-session
autostart=true
autorestart=false
priority=20
startsecs=5
startretries=3
user=root
environment=HOME="/root",USER="root",DISPLAY=":1",DBUS_SESSION_BUS_ADDRESS="autolaunch:"
stdout_logfile=/var/log/supervisor/xfce.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/xfce-err.log
stderr_logfile_maxbytes=5MB
