# modules/vscode/supervisord.conf
[program:vscode]
; code-server — VS Code in the browser on port 8200.
; No authentication for local/sandboxed use.
; Accessible via nginx at /vscode/
command=/usr/bin/code-server
    --bind-addr=0.0.0.0:8200
    --auth=none
    --user-data-dir=/root/.local/share/code-server
    --extensions-dir=/root/.local/share/code-server/extensions
    /root
autostart=true
autorestart=true
priority=50
startsecs=5
startretries=3
user=root
environment=HOME="/root",USER="root"
stdout_logfile=/var/log/supervisor/vscode.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/vscode-err.log
stderr_logfile_maxbytes=5MB
