# modules/sms/supervisord.conf
[program:sms-webhook]
; uvicorn running the FastAPI SMS webhook on port 8100.
; Requires ANTHROPIC_API_KEY and WHATSAPP_RECIPIENT to be set in the environment.
command=/usr/bin/python3 -m uvicorn sms_webhook:app
    --host 0.0.0.0
    --port 8100
    --workers 1
    --log-level info
    --no-access-log
directory=/sandbox/modules/sms
autostart=true
autorestart=true
priority=50
startsecs=5
startretries=3
user=root
environment=
    HOME="/root",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    WHATSAPP_RECIPIENT="%(ENV_WHATSAPP_RECIPIENT)s",
    WHATSAPP_WEBHOOK_URL="%(ENV_WHATSAPP_WEBHOOK_URL)s"
stdout_logfile=/var/log/supervisor/sms-webhook.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/sms-webhook-err.log
stderr_logfile_maxbytes=5MB
