# modules/teams/supervisord.conf
[program:teams]
; Opens Microsoft Teams as a PWA in Chromium using a dedicated profile.
; Uses the shared display :1 so it appears in the noVNC desktop.
; Profile is isolated at /root/.config/teams-profile to avoid conflicting
; with the browser module's chromium-sandbox profile.
command=/usr/bin/chromium-browser
    --no-sandbox
    --disable-dev-shm-usage
    --disable-gpu
    --disable-software-rasterizer
    --user-data-dir=/root/.config/teams-profile
    --profile-directory=Default
    --app=https://teams.microsoft.com
    --window-size=1280,800
    --window-position=0,0
    --no-first-run
    --no-default-browser-check
    --disable-background-networking
    --disable-extensions
    --disable-sync
    --disable-translate
    --password-store=basic
    --enable-features=WebRTC
    --use-fake-ui-for-media-stream
autostart=true
autorestart=false
priority=60
startsecs=5
startretries=2
user=root
environment=DISPLAY=":1",HOME="/root",DBUS_SESSION_BUS_ADDRESS="autolaunch:"
stdout_logfile=/var/log/supervisor/teams.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/teams-err.log
stderr_logfile_maxbytes=5MB
